# rootfs/etc/services.d/caliver/run
#!/usr/bin/with-contenv bash

HOME=~app
export HOME

exec s6-setuidgid app /app/calibre/calibre-server --${AUTH_STATUS} --userdb ${CONF_DIR}/files/${USERSQL} ${ADD_LIB}